# config/files/usr/etc/systemd/system/rpm-ostreed-automatic.timer.d/override.conf
[Timer]
RandomizedDelaySec=10m
OnCalendar=*-*-* 17:45:00 UTC
Persistent=true
