# config/files/usr/etc/systemd/user/distrobox-auto-update.service
[Unit]
Description=Automatically update all distroboxes
[Service]
Type=simple
ExecStart=distrobox-upgrade --all
StandardOutput=null
