# config/files/usr/etc/systemd/user/distrobox-auto-update.timer
[Timer]
RandomizedDelaySec=10m
OnCalendar=*-*-* 05:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
